# deploy/systemd/daphne.service.example
# /etc/systemd/system/daphne.service
# Replace User and every path segment: clone folder may be chat_bot, chatbot_emotion_detector, etc.

[Unit]
Description=Daphne ASGI (Emotion Chat)
After=network.target

[Service]
Type=simple
User=REPLACE_USER
Group=REPLACE_USER
WorkingDirectory=/home/REPLACE_USER/chat_bot/emotion_chat
Environment=PATH=/home/REPLACE_USER/chat_bot/.venv/bin
ExecStart=/home/REPLACE_USER/chat_bot/.venv/bin/daphne \
    -b 127.0.0.1 -p 8000 emotion_chat.asgi:application
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
